# cliphist/.config/systemd/user/cliphist-wipe.timer
[Unit]
Description=Weekly cliphist history wipe

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
