# cliphist/.config/systemd/user/cliphist-wipe.service
[Unit]
Description=Wipe cliphist clipboard history

[Service]
Type=oneshot
ExecStart=/run/current-system/sw/bin/cliphist wipe
